# setup/thoth-firstboot.service
[Unit]
Description=Thoth First Boot Setup
After=network-online.target
Wants=network-online.target
Before=thoth-web.service
ConditionPathExists=!/etc/thoth-first-boot-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 5
ExecStart=/bin/bash /home/pi/thoth/setup/first-boot.sh
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
